# systemd/zerolayer.timer
[Unit]
Description=Rebuild your image to keep you up to date.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
